# timeupdate.service/updatetime.timer
[Unit]
Description=Run clockupdate every 1min

[Timer]
OnUnitActiveSec=1min
OnBootSec=10min
Unit=updatetime

[Install]
WantedBy=timers.target
